# systemd/rocksky-scrobble-abuse-sweep.service
[Unit]
Description=Rocksky Scrobble Abuse Sweep Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/github/rocksky/apps/scrobble-abuse-sweep
ExecStart=/bin/bash -ic 'exec doppler run -- deno task start'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rocksky-scrobble-abuse-sweep

# Environment
Environment=HOME=/root
Environment=USER=root

[Install]
WantedBy=multi-user.target
