# systemd/rocksky-jetstream.service
[Unit]
Description=Rocksky Jetstream Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/github/rocksky
ExecStart=/bin/bash -ic 'exec doppler run -- ./target/release/rockskyd jetstream'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rocksky-jetstream

# Environment
Environment=HOME=/root
Environment=USER=root

[Install]
WantedBy=multi-user.target
